# ast.timer
[Unit]
Description=Runs AST Hotline sending every 2 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
